# cheddar-worker.service
[Unit]
Description=Cheddar Logic Worker Scheduler
Documentation=file:///opt/cheddar-logic/PRODUCTION_CUTOVER.md
After=network.target
Wants=cheddar-logic-db.service

[Service]
Type=simple
User=babycheeses11
WorkingDirectory=/opt/cheddar-logic
Environment="NODE_ENV=production"
Environment="DATABASE_PATH=/opt/data/cheddar.db"
Environment="CHEDDAR_DB_PATH=/opt/data/cheddar.db"
Environment="CHEDDAR_DATA_DIR=/opt/data"
EnvironmentFile=/opt/cheddar-logic/.env.production

ExecStart=/usr/bin/npm --prefix /opt/cheddar-logic/apps/worker run scheduler

# Restart policy
Restart=on-failure
RestartSec=30s
StartLimitInterval=300s
StartLimitBurst=5

# Shutdown handling
TimeoutStopSec=10s
KillMode=mixed
KillSignal=SIGTERM

# Resource limits
MemoryLimit=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cheddar-worker

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/data /opt/cheddar-logic/packages/data

[Install]
WantedBy=multi-user.target
