# deploy/systemd/cheddar-web.service
[Unit]
Description=Cheddar Logic Web App (Next.js)
After=network.target

[Service]
Type=simple
User=babycheeses11
WorkingDirectory=/opt/cheddar-logic/web
Environment="NODE_ENV=production"
Environment="PORT=3000"
EnvironmentFile=/opt/cheddar-logic/.env.production

ExecStart=/usr/bin/node /opt/cheddar-logic/web/node_modules/.bin/next start --port 3000

# Restart policy
Restart=on-failure
RestartSec=10s
# Web process is read-only (no DB writes). Kill fast so deploys don't hit the
# 90-second systemd default timeout and get SIGKILLed mid-flight.
TimeoutStopSec=15s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cheddar-web

# Security
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/data /opt/cheddar-logic/packages/data /opt/cheddar-logic/web/.next /tmp

[Install]
WantedBy=multi-user.target
