# deploy/systemd/cheddar-web.service
[Unit]
Description=Cheddar Logic Web App (Next.js)
After=network.target

[Service]
Type=simple
User=babycheeses11
WorkingDirectory=/opt/cheddar-logic/web
Environment="NODE_ENV=production"
Environment="PORT=3000"
EnvironmentFile=/opt/cheddar-logic/.env.production

ExecStart=/usr/bin/node /opt/cheddar-logic/web/node_modules/.bin/next start --port 3000
# Restart cloudflared AFTER Next.js is ready on port 3000 (not just after process fork).
#
# WHY: ExecStartPost fires at process start time (Type=simple). Next.js takes
# 5-15s to cold-start. Restarting cloudflared before the origin is ready means
# the tunnel is live but 127.0.0.1:3000 returns ECONNREFUSED → 502 for users.
#
# This script polls port 3000 until it responds, then restarts cloudflared.
# The + prefix runs as root (required to call systemctl).
ExecStartPost=+/opt/cheddar-logic/deploy/restart-cloudflared-when-ready.sh

# Restart policy
Restart=on-failure
RestartSec=10s
# Web process is read-only (no DB writes). Allow enough time for Next.js to
# drain in-flight requests before systemd falls back to SIGKILL.
# 15s was too short — production loads caused SIGKILL on every deploy restart.
TimeoutStopSec=30s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cheddar-web

# Security
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/data /opt/cheddar-logic/packages/data /opt/cheddar-logic/web/.next /tmp

[Install]
WantedBy=multi-user.target
