# deploy/systemd/cheddar-web.service
[Unit]
Description=Cheddar Logic Web App (Next.js)
After=network.target

[Service]
Type=simple
User=babycheeses11
WorkingDirectory=/opt/cheddar-logic/web
Environment="NODE_ENV=production"
Environment="PORT=3000"
EnvironmentFile=/opt/cheddar-logic/.env.production

ExecStart=/usr/bin/node /opt/cheddar-logic/web/node_modules/.bin/next start --port 3000
# ExecStartPost cloudflared restart intentionally removed.
#
# WHY IT WAS HERE (wrong): the prior script restarted cloudflared after
# Next.js became ready, on the theory that cloudflared needed a fresh
# connection to the newly-started origin. This is incorrect and harmful:
# cloudflared maintains its own tunnel and connection pool; when the origin
# recovers it resumes proxying automatically. Restarting cloudflared on every
# web-process lifecycle event added a SECOND 5-30s 502 window on top of the
# unavoidable RestartSec gap, doubling every outage. Cloudflare returns 502
# naturally while the origin is down; no cloudflared restart is needed.

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitInterval=300s
StartLimitBurst=5
# Web process is read-only (no DB writes). Allow enough time for Next.js to
# drain in-flight requests before systemd falls back to SIGKILL.
TimeoutStopSec=30s
# NODE_OPTIONS: tell V8 to GC aggressively when heap approaches 614 MB
# (80% of the MemoryLimit ceiling), so Next.js GCs before systemd kills the
# process. Without this, V8 may grow past the OOM trigger before running GC.
Environment="NODE_OPTIONS=--max-old-space-size=614"
MemoryLimit=768M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cheddar-web

# Security
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/data /opt/cheddar-logic/packages/data /opt/cheddar-logic/web/.next /tmp

[Install]
WantedBy=multi-user.target
